# progname.conf
[program:progname]
directory=/path/to/code/dir
command=/path/to/code/dir/_main.py
user=username
autostart=true
autorestart=true